# a-bcached.service
[Unit]
Description=a-bcache HDD spin management daemon
Documentation=https://github.com/asymysh/a-bcache
After=local-fs.target
Requires=local-fs.target

[Service]
Type=forking
ExecStart=/usr/sbin/a-bcached \
    --spindown-delay 300 \
    --min-writeback 64 \
    --poll-interval 5 \
    --max-spin-cycles 4
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=10

# Security hardening
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/sys/block /sys/fs/bcache
NoNewPrivileges=false
# Needs CAP_SYS_RAWIO for hdparm, CAP_DAC_OVERRIDE for sysfs writes
AmbientCapabilities=CAP_SYS_RAWIO CAP_DAC_OVERRIDE

[Install]
WantedBy=multi-user.target
