# a-bcached.service
[Unit]
Description=a-bcache HDD spin management daemon
Documentation=https://github.com/asymysh/a-bcache
After=local-fs.target
Wants=local-fs.target

[Service]
# Type=simple with --foreground avoids the PIDFile tracking problems
# associated with Type=forking.
Type=simple
ExecStart=/usr/sbin/a-bcached \
    --foreground \
    --spindown-delay 300 \
    --min-writeback 64 \
    --poll-interval 5 \
    --max-spin-cycles 4

Restart=on-failure
RestartSec=10

# Runs as root because hdparm needs CAP_SYS_RAWIO and writes to bcache
# sysfs control files. Note: we cannot set ProtectKernelTunables=true
# because that makes /sys read-only, which would break sysfs writes.
PrivateTmp=true
ProtectHome=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
